# systemd/postmark-tti.timer
[Unit]
Description=Timer to fetch Postmark Time to Inbox data

[Timer]
OnBootSec=1min
OnUnitActiveSec=1hour

[Install]
WantedBy=default.target
